# templates/default/x11vnc.conf.erb
start on login-session-start
script
  start-stop-daemon --start --exec /usr/bin/x11vnc -- -bg -display <%= @display %> -o <%= @log_file %> <%= @executable_args %> <%= @password_opt %> -auth /var/run/lightdm/root/:0
end script
