# templates/default/x11vnc.init.d.erb
#!/bin/sh
# $Id: Generated by chef for node: <%= node['hostname'] %>
# ----------------------------------------------------------------
# NOTE: This file is controlled by chef templates!
# Do not edit or change this file but change the following:
# template file : <%= @template_file %>
# recipe file : <%= @recipe_file %>
# ----------------------------------------------------------------
#
# x11vnc VNC access to X11
#
# chkconfig:   - 70 30
# description: Gives access to a running X11 server via VNC

# Source function library.
. /etc/rc.d/init.d/functions

exec="/usr/bin/x11vnc"
args="-bg -display <%= @display %> -o <%= @log_file %> <%= @executable_args %> <%= @password_opt %>"
prog="x11vnc"

[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

lockfile=/var/lock/subsys/$prog

start() {
    [ -x $exec ] || exit 5
    echo -n $"Starting $prog: "
    daemon $exec $args
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    return $retval
}

stop() {
    echo -n $"Stopping $prog: "
    killproc $prog
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    stop
    start
}

reload() {
    restart
}

force_reload() {
    restart
}

rh_status() {
    # run checks to determine if the service is running or use generic status
    status $prog
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}

case "$1" in
    start)
        rh_status_q && exit 0
        $1
        ;;
    stop)
        rh_status_q || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        rh_status_q || exit 7
        $1
        ;;
    force-reload)
        force_reload
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q || exit 0
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart|reload|force-reload}"
        exit 2
esac
exit $?
